# supervisord.conf
[supervisord]
nodaemon=true

[program:nodepki]
directory=/root/nodepki/
command=node server.js

[program:nodepki_webclient]
directory=/root/nodepki-webclient/
command=node app.js
