# supervisord.conf
[supervisord]
nodaemon=true

[program:nodepki]
directory=/root/nodepki/
command=nodejs server.js

[program:nodepki_webclient]
directory=/root/nodepki-webclient/
command=nodejs app.js
