# production/supervisord.conf
[supervisord]
nodaemon=true

[program:nodepki]
directory=/opt/nodepki
command=node server.js

[program:nodepki_webclient]
directory=/opt/nodepki-webclient/
command=node app.js
